# com.happyquotes.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.happyquotes.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__HAPPYQUOTES_PATH__/happyquotes.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HAPPYQUOTES_PATH__</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/happyquotes.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/happyquotes.log</string>
</dict>
</plist>
